# roles/partynest/templates/partynest.service
[Unit]
After=network.target
Description={{ partynest__general_descr }}

[Service]
ExecStart=/bin/true
ExecReload=/bin/true
RemainAfterExit=yes
Type=oneshot
WorkingDirectory=/opt/partynest/current

[Install]
WantedBy=multi-user.target
