# roles/partynest/templates/partynest-restart.path
[Path]
PathModified=/opt/partynest/current

[Install]
WantedBy=multi-user.target
