# roles/partynest/templates/partynest-restart.path
[Path]
PathModified=/opt/partynest/revisions.log

[Install]
WantedBy=multi-user.target
